# package/usr/lib/systemd/system/healthchecks@.timer
[Unit]
Description=Send heartbeat every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=true

[Install]
WantedBy=timers.target